# supervisor-app.conf
[supervisord]
nodaemon=true

[supervisorctl]

[program:app-uwsgi]
command = uwsgi --ini uwsgi.ini

[program:nginx-app]
command = nginx -c {ABSOLUTE_PROJECT_PATH}nginx-app.conf -g "daemon off;"
